# data/rm-pad@.service
[Unit]
Description=rm-pad input forwarding for reMarkable tablet
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rm-pad
Restart=on-failure
RestartSec=5
Environment="RMPAD_HOST=10.11.99.1"

[Install]
WantedBy=default.target
